# systemd-socket-activation/h5g-u1.socket
[Socket]
ListenStream=/tmp/u1-systemd.socket

[Install]
WantedBy=sockets.target
